# test/phpsq.conf
[program:phpsq_test]
command=/var/www/phpsq/test/phpsqworker.php queue=test
directory=.
user=nobody
process_name=%(program_name)s_%(process_num)03d
numprocs=3
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s_%(process_num)03d_out
stderr_logfile=/var/log/supervisor/%(program_name)s_%(process_num)03d_err